# rootfs/etc/s6-overlay/s6-rc.d/tor/run
#!/command/with-contenv sh
exec /usr/bin/tor -f /etc/tor/torrc
