# rootfs/etc/s6-overlay/s6-rc.d/torproxy-config/up
#!/command/execlineb -P
with-contenv

/usr/bin/torproxy.sh

if [ -e /var/lib/tor/hidden_service/hostname ]; then
    echo ""
    echo -n "Hidden service hostname: "
    cat /var/lib/tor/hidden_service/hostname
fi
